# ansible/files/etc/systemd/system/nginx-watcher.path
[Path]
Unit=nginx-reloader.service
PathChanged=/etc/pve/local

[Install]
WantedBy=multi-user.target